# server/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/home/cypher/nowallet/endpoint.sock

[Install]
WantedBy=sockets.target
